# pipservice-cap.service
[Unit]
Description=ASF Contingent Approval Service

[Service]
Type=exec
WorkingDirectory=/opt/cap
ExecStart=/usr/local/bin/uv run cap-backend --config /opt/cap/config.yaml
Restart=always
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target